# detections/T1053.005_scheduled_task.yaml
title: Scheduled Task Created
id: a1b2c3d4-0003-0000-0000-000000000003
status: experimental
description: Detects creation of a scheduled task (Security 4698), used for persistence/execution.
references:
  - https://attack.mitre.org/techniques/T1053/005/
tags:
  - attack.persistence
  - attack.t1053.005
logsource:
  product: windows
  service: security
detection:
  selection:
    EventID: 4698
  condition: selection
falsepositives:
  - Legitimate software/admin creating scheduled tasks
level: medium